# misc/service/log/run
#!/bin/bash

exec setuidgid ata multilog ./main
